# graham/deploy/graham-weekly.timer
[Unit]
Description=Run GRAHAM weekly on Sunday at 02:00 UTC

[Timer]
OnCalendar=Sun 02:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
