# scripts/com.atlas.autonomous.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.atlas.autonomous</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-m</string>
        <string>agents.autonomous_loop</string>
        <string>--once</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/chrisworsey/Desktop/atlas</string>
    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run at 4:30pm ET (21:30 UTC) on weekdays -->
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/chrisworsey/Desktop/atlas/logs/autonomous.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/chrisworsey/Desktop/atlas/logs/autonomous_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
